# mirror/CRAN/QA/BDR/blackswan/crontab
5 4-23/6 * * * /data/blackswan/ripley/R/packages/check-gcc > /data/blackswan/ripley/R/packages/check-gcc.log 2>&1
