# mirror/CRAN/QA/BDR/tawny/crontab
5 2 * * * /Users/ripley/R/packages/check-devel > /Users/ripley/R/packages/check-devel.log 2>&1
