# mirror/CRAN/QA/BDR/M1mac/crontab
5 1 * * * /Users/ripley/R/packages/tl-update > /Users/ripley/R/packages/tl-update.log  2>&1
