# mirror/CRAN/QA/BDR/Solaris/Sparc/crontab
25 3,9,15,21 * * * /home/ripley/R/packages/check-Sparc > /home/ripley/R/packages/check-Sparc.log 2>&1
